# app.service
[Unit]
Description=Domain Liveness Check Flask App
After=network.target
Wants=network-online.target

[Service]
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/domain-monitor-system

# Run the app using the virtual environment Python
ExecStart=/home/ubuntu/domain-monitor-system/_venv_/bin/python3 /home/ubuntu/domain-monitor-system/app.py

# Environment variables
Environment="PATH=/home/ubuntu/domain-monitor-system/_venv_/bin"
Environment="FLASK_ENV=production"
Environment="PYTHONUNBUFFERED=1"

# Log to journal
StandardOutput=journal
StandardError=journal

# Restart policy
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
