# app.service
#!/bin/bash
set -e  # exit on any error

# -----------------------------
# Function to print status
# -----------------------------
function check_result() {
    if [ $? -ne 0 ]; then
        echo "[ERROR] $1 failed. Exiting."
        exit 1
    fi
    echo "[SUCCESS] $1 succeeded."
}

echo "[INFO] Starting deployment of Domain Monitor System..."

# -----------------------------
# Update system packages
# -----------------------------
echo "[INFO] Updating repositories..."
sudo apt update -y
check_result "Updating repositories"

echo "[INFO] Installing dependencies..."
sudo apt install -y git python3-venv python3-pip
check_result "Installing dependencies"

# -----------------------------
# Clone or update repo
# -----------------------------
APP_DIR="$HOME/domain-monitor-system"

if [ -d "$APP_DIR" ]; then
    echo "[INFO] Repo exists. Pulling latest changes..."
    cd "$APP_DIR"
    git reset --hard
    git pull
else
    echo "[INFO] Cloning repository..."
    git clone https://github.com/Tomer-ui/domain-monitor-system.git "$APP_DIR"
    cd "$APP_DIR"
fi
check_result "Repo ready"

# -----------------------------
# Create virtual environment
# -----------------------------
if [ ! -d "_venv_" ]; then
    echo "[INFO] Creating Python virtual environment..."
    python3 -m venv _venv_
    check_result "Virtual environment created"
else
    echo "[INFO] Virtual environment already exists."
fi

# Activate venv and install requirements
echo "[INFO] Installing Python dependencies..."
source _venv_/bin/activate
pip install --upgrade pip
pip install -r requirements.txt
check_result "Python dependencies installed"

# -----------------------------
# Systemd service
# -----------------------------
echo "[INFO] Setting up systemd service..."
SERVICE_FILE="/etc/systemd/system/app.service"

# Write the service file (idempotent)
sudo bash -c "cat > $SERVICE_FILE" <<EOL
[Unit]
Description=Domain Liveness Check Flask App
After=network.target
Wants=network-online.target

[Service]
User=ubuntu
Group=ubuntu
WorkingDirectory=$APP_DIR
ExecStart=$APP_DIR/_venv_/bin/python3 $APP_DIR/app.py
Environment=\"PATH=$APP_DIR/_venv_/bin\"
Environment=\"FLASK_ENV=production\"
Environment=\"PYTHONUNBUFFERED=1\"
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
EOL

sudo systemctl daemon-reload
sudo systemctl enable app.service
sudo systemctl restart app.service
check_result "Systemd service installed and started"

echo "[INFO] Deployment complete! Flask app is running on port 8080."
